# workspaces/dating/setup/com.shakstzy.quantum-tinder-pull.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Pulls matches + thread snapshots, then runs decide. Twice a day, never within 90 min of a swipe session. -->
<plist version="1.0">
<dict>
  <key>Label</key><string>com.shakstzy.quantum-tinder-pull</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>cd /Users/shakstzy/QUANTUM/workspaces/tinder/bot &amp;&amp; sleep $((RANDOM % 1800)) &amp;&amp; node scripts/pull.mjs &amp;&amp; node scripts/decide.mjs &gt;&gt; ~/Library/Logs/quantum-tinder.log 2&gt;&amp;1</string>
  </array>
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Hour</key><integer>15</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>23</integer><key>Minute</key><integer>0</integer></dict>
  </array>
  <key>StandardOutPath</key><string>/Users/shakstzy/Library/Logs/quantum-tinder.stdout.log</string>
  <key>StandardErrorPath</key><string>/Users/shakstzy/Library/Logs/quantum-tinder.stderr.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key><string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
