# _shell/plists/.unpatched-backup/com.adithya.ultron.dating-bumble-send.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- 15 fires/day with non-grid minute marks + 0-30min jitter on top.
     Effective fire window per slot: ~60min wide, irregular spacing.
     Defeats "obvious cron loop" pattern detection. Hand-written; orphan to schedule skill. -->
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>HOME</key><string>/Users/shakstzy</string>
		<key>PATH</key><string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
		<key>ULTRON_ROOT</key><string>/Users/shakstzy/ULTRON</string>
	</dict>
	<key>Label</key><string>com.adithya.ultron.dating-bumble-send</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>-c</string>
		<string>flock -n /tmp/com.adithya.ultron.dating-bumble-send.lock bash -c 'cd /Users/shakstzy/ULTRON/workspaces/dating/bot-bumble &amp;&amp; sleep $((RANDOM % 1800)) &amp;&amp; node scripts/send.mjs'</string>
	</array>
	<key>RunAtLoad</key><false/>
	<key>StandardErrorPath</key><string>/Users/shakstzy/ULTRON/_logs/com.adithya.ultron.dating-bumble-send.err.log</string>
	<key>StandardOutPath</key><string>/Users/shakstzy/ULTRON/_logs/com.adithya.ultron.dating-bumble-send.out.log</string>
	<key>StartCalendarInterval</key>
	<array>
		<dict><key>Hour</key><integer>9</integer><key>Minute</key><integer>8</integer></dict>
		<dict><key>Hour</key><integer>10</integer><key>Minute</key><integer>23</integer></dict>
		<dict><key>Hour</key><integer>11</integer><key>Minute</key><integer>46</integer></dict>
		<dict><key>Hour</key><integer>12</integer><key>Minute</key><integer>34</integer></dict>
		<dict><key>Hour</key><integer>13</integer><key>Minute</key><integer>17</integer></dict>
		<dict><key>Hour</key><integer>14</integer><key>Minute</key><integer>42</integer></dict>
		<dict><key>Hour</key><integer>15</integer><key>Minute</key><integer>28</integer></dict>
		<dict><key>Hour</key><integer>16</integer><key>Minute</key><integer>51</integer></dict>
		<dict><key>Hour</key><integer>17</integer><key>Minute</key><integer>33</integer></dict>
		<dict><key>Hour</key><integer>18</integer><key>Minute</key><integer>19</integer></dict>
		<dict><key>Hour</key><integer>19</integer><key>Minute</key><integer>42</integer></dict>
		<dict><key>Hour</key><integer>20</integer><key>Minute</key><integer>14</integer></dict>
		<dict><key>Hour</key><integer>21</integer><key>Minute</key><integer>38</integer></dict>
		<dict><key>Hour</key><integer>22</integer><key>Minute</key><integer>9</integer></dict>
		<dict><key>Hour</key><integer>22</integer><key>Minute</key><integer>51</integer></dict>
	</array>
	<key>WorkingDirectory</key><string>/Users/shakstzy/ULTRON/workspaces/dating/bot-bumble</string>
</dict>
</plist>
